# etc/init.d/servofer
#!/bin/sh /etc/rc.common

START=99

boot() {
	node /root/servoDeFer/ &
}
start() { 
	node /root/servoDeFer/ &
	}
stop() {
	echo "servoDeFer stopped"
	}