# etc/init.d/sketchpad
#!/bin/sh /etc/rc.common

START=99

boot() {
	node /root/sketchpad/ &
}
start() { 
	node /root/sketchpad/ &
	}
stop() {
	echo "sketchpad stopped"
	}